# opt/profile.d/pulseaudio.csh
setenv SDL_AUDIODRIVER pulse
setenv AUDIODRIVER pulseaudio
setenv QEMU_AUDIO_DRV pa
setenv ALSA_CONFIG_PATH /etc/alsa-pulse.conf
