# opt/profile.d/pulseaudio.sh
alias aumix='padsp aumix'
alias sox='padsp sox'
export SDL_AUDIODRIVER=pulse
alias timidity='timidity -Oe'
export AUDIODRIVER=pulseaudio
export QEMU_AUDIO_DRV=pa
export ALSA_CONFIG_PATH=/etc/alsa-pulse.conf
